# tests/overlay/etc/systemd/system/net-setup.service
[Service]
Type=oneshot
ExecStart=/usr/local/bin/net-setup.sh
[Install]
WantedBy=sshd.service