# tests/overlay/etc/systemd/system/scout-setup.service
[Service]
Type=oneshot
ExecStart=/usr/local/bin/scout-setup.sh

[Install]
WantedBy=sysinit.target
